# deploy/supervisor.conf
[program:gunicorn]
command=/restcar/venv/bin/gunicorn manage:app -c /restcar/gunicorn.conf.py
directory=/restcar/venv/
environment=FLASK_ENV="live"
